# resources/docs/platform/nginx-go/docker/config/crontab
* * * * * /usr/bin/docker exec myapp_container /usr/local/bin/app schedule >> /var/log/myapp-schedule.log 2>&1